# services/bridge/run
#!/bin/bash

set -euxo pipefail
cd /opt/slackbridge
exec /opt/slackbridge/venv/bin/python -m slackbridge
